# mysql-wp/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:zabbix-agent2]
command=zabbix_agent2 -c /etc/zabbix/zabbix_agent2.conf
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
